# multiwan-nft/files/etc/uci-defaults/multiwan-nft-migrate-flush_conntrack
#!/bin/sh

exec 2>/dev/null

. /lib/functions.sh

command -v uci >/dev/null 2>&1 || exit 0
[ -f /etc/config/multiwan-nft ] || exit 0

changed=0

multiwan_nft_migrate_flush_conntrack() {
	local iface="$1"

	config_get value "${iface}" flush_conntrack
	case $value in
		always)
			uci -q delete "multiwan-nft.${iface}.flush_conntrack" 2>/dev/null || true
			uci -q add_list "multiwan-nft.${iface}.flush_conntrack=ifup" 2>/dev/null || true
			uci -q add_list "multiwan-nft.${iface}.flush_conntrack=ifdown" 2>/dev/null || true
			changed=1
			;;
		never)
			uci -q delete "multiwan-nft.${iface}.flush_conntrack" 2>/dev/null || true
			changed=1
			;;
	esac
}

config_load 'multiwan-nft' 2>/dev/null || exit 0
config_foreach multiwan_nft_migrate_flush_conntrack interface

[ "$changed" = "1" ] && uci -q commit multiwan-nft 2>/dev/null || true

exit 0
